# deploy/livekit-backend.service
[Unit]
Description=Gunicorn for LiveKit demo Django backend
After=network.target docker.service

[Service]
User=opc
WorkingDirectory=/opt/livekit-demo/backend
EnvironmentFile=/opt/livekit-demo/backend/.env
ExecStart=/opt/livekit-demo/backend/venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:8000 --workers 3
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
